# deploy/gcp/trading-bot.service
[Unit]
Description=Trading Bot Telegram Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=__PYTHON_BIN__ __APP_DIR__/start_telegram_bot.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
